# fancurve.service
[Unit]
Description=fancurve service.

[Service]
Type=simple
ExecStart=/usr/local/bin/fancurve

[Install]
WantedBy=multi-user.target